# ansible/files/mlnx.service
[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/opt/mlnx-update.sh

[Install]
WantedBy=multi-user.target
